# systemd/usbguard-web.service
[Unit]
Description=USBGuard Approval Manager - Web Interface
Documentation=https://github.com/usbguard/usbguard-manager
Requires=usbguard.service
After=network.target usbguard.service
Wants=usbguard-ttl-reaper.timer
# ⬇️ הוסף את השורה הבאה - מונע בעיות אם הרשת עולה מאוחר
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=usbadmins
WorkingDirectory=/etc/usbguard/web
ExecStart=/etc/usbguard/web/start-web.sh
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectHome=true
ProtectSystem=strict
# ⬆️ שיניתי מ-full ל-strict (מאפשר ReadWritePaths לעבוד טוב יותר)
PrivateTmp=true
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# התיקון שלך - תקין לחלוטין
ReadWritePaths=/etc/usbguard /var/lib/usbguard-manager /var/log /var/tmp

# ⬇️ הוסף timeout כדי למנוע hang בהשבתה
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target